# crontabs
5       0       *       *       *       /bin/sh /opt/generate
